# pkg/arch/leosac.install
#!/bin/sh

post_install()
{
    echo "----------------------------------------------------------------"
    echo "To enable leosac daemon run these commands as root:"
    echo "systemctl enable leosacd"
    echo "systemctl start leosacd"
    echo ""
    echo "More info here:"
    echo "https://github.com/islog/leosac/wiki"
    echo "----------------------------------------------------------------"
}

post_upgrade()
{
    echo "----------------------------------------------------------------"
    echo "To restart leosac daemon run these commands as root:"
    echo "systemctl daemon-reload"
    echo "systemctl restart leosacd"
    echo ""
    echo "More info here:"
    echo "https://github.com/islog/leosac/wiki"
    echo "----------------------------------------------------------------"
}
